# admin-web/deploy/divehub-admin-web.service.example
[Unit]
Description=DiveHub admin-web (Next.js production)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/HUB/admin-web
Environment=NODE_ENV=production
Environment=PORT=3003
# nvm: замените на полный путь, например ExecStart=/root/.nvm/versions/node/v22.0.0/bin/npm run start
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
# При запуске не от root раскомментируйте и выровняйте владельца каталога admin-web:
# User=www-data
# Group=www-data

[Install]
WantedBy=multi-user.target
